# inputleap/inputleap-autostart.service
[Unit]
Description=InputLeap Client Auto-start
Before=display-manager.service
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
User=gdm
Group=gdm
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/120
Environment=HOME=/var/lib/gdm3
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/flatpak run --user io.github.input_leap.input-leap
Restart=always
RestartSec=3

[Install]
WantedBy=graphical.target